# roles/jetson_aip_prereqs/templates/aip-usb-stage.path.j2
[Unit]
Description=Watch for AIP_DEPLOY USB volume

[Path]
PathExistsGlob=/media/*/{{ aip_usb_label }}

[Install]
WantedBy=multi-user.target
