# roles/jetson_app_stack/files/docker-compd/docker-compose-pos.service
[Unit]
Description=Start POS Application After Xorg & PulseAudio
Wants=docker.service xorg-pulseaudio.service
Requires=docker.service xorg-pulseaudio.service
After=docker.service xorg-pulseaudio.service pulseaudio.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/checkout/eco_boutique_deploy

# Remove orphaned containers but keep data
ExecStartPre=/usr/bin/docker compose down --remove-orphans

# Start without rebuilding each time
ExecStart=/usr/bin/docker compose up -d

# Stop gracefully
ExecStop=/usr/bin/docker compose down

# Restart only on failure (prevents infinite loops)
Restart=on-failure
RestartSec=10s

User=checkout
Group=checkout

[Install]
WantedBy=multi-user.target
